# misc/dde-clipboard-daemon.service.in
[Unit]
Description=dde-clipboard-daemon
Requisite=dde-session-pre.target
After=dde-session-pre.target

Requisite=dde-session-initialized.target
PartOf=dde-session-initialized.target
Before=dde-session-initialized.target
StartLimitBurst=3

[Service]
Type=dbus
BusName=org.deepin.dde.daemon.Clipboard1
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/dde-clipboard-daemon
Slice=session.slice
Restart=on-failure
RestartSec=3s
