# misc/dde-clipboard.service.in
[Unit]
Description=Deepin Clipboard

[Service]
Type=dbus
BusName=org.deepin.dde.Clipboard1
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/dde-clipboard
Slice=session.slice
